# services/exif-orientation-service/dist/exif-orientation-service.conf
description "EXIF Orientation Service"
author "hoccer"

env HOME=/home/talk/exif-orientation-service/current
setuid talk
setgid talk

start on runlevel [2345]
stop on runlevel [06]

script
  cd $HOME
  /usr/local/rvm/bin/rvm 2.1.1@exif-orientation-service do bundle exec ruby service.rb -e production
  emit exif_orientation_service_running
end script
