# dist/hoccer-receiver.conf
description "Hoccer Receiver Service"
author "hoccer"

env HOME=/home/receiver/hoccer-receiver/current
setuid receiver
setgid receiver

start on runlevel [2345]
stop on runlevel [06]

script
  cd $HOME
  java -jar hoccer-receiver.jar -c config/receiver.conf 2>&1 >> log/hoccer-receiver.log
  emit receiver_running
end script
