# dist/talk-webclient-backend.conf
description "Hoccer XO WebClient Backend Service"
author "hoccer"

env HOME=/home/talk/webclient-backend/current
setuid talk
setgid talk

start on runlevel [2345]
stop on runlevel [06]

script
  cd $HOME
  java -jar talk-webclient-backend.jar -c backend.conf 2>&1 >> log/webclient-backend.log
  emit webclient_backend_running
end script
